# web_controller/game-controller.service
[Unit]
Description=Game Controller WebSocket Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/c/dev/VTHacks13/web_controller
ExecStart=/usr/bin/python3 /mnt/c/dev/VTHacks13/web_controller/server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
